# fcos/compose-up@.path
[Unit]
Description=Run Compose Images for %I

[Path]
PathModified=/etc/compose/build/%i/compose.yml
Unit=compose-up@%i.service

[Install]
WantedBy=multi-user.target
