# fcos/compose-build@.path
[Unit]
Description=Build Compose Images for %I

[Path]
PathChanged==/etc/compose/build/%i
Unit=compose-build@%i.service

[Install]
WantedBy=multi-user.target
