# container/knot@.container
[Unit]
Description=Knot DNS Server (%i)
After=container-build@%N.service container-secrets.service podman.service
Wants=container-build@%N.service container-secrets.service podman.service

[Service]
# environment loaded here is available in systemd-quadlet scope
EnvironmentFile=-/etc/containers/environment/%N-%i-systemd.env
Restart=on-failure

[Container]
Image=localhost/%N:latest
# If an image is updated in local storage, Podman restarts the systemd unit
AutoUpdate=local

# Use instance name in container name
ContainerName=%N-%i
# Instance-specific volume name
Volume=%N-%i-data:/var/lib/knot:rw
# Mount instance-specific config and zone files
Volume=/etc/local/%N/%i:/etc/knot:ro
# Run as knot user inside the container
User=knot
Group=knot
# Expose port 53 (default), will later be overridden by instance
PublishPort=53:53/udp
PublishPort=53:53/tcp

# Example: Run a command after the container starts (adjust as needed).
# This can be useful for verifying the config or initializing something.
#Exec=/usr/sbin/knotc -c /etc/knot/knot.conf conf-check

[Install]
WantedBy=multi-user.target
